# supervisord.conf
[supervisord]
nodaemon=true

[program:squid]
command         = /usr/sbin/squid3 -f /etc/squid/squid.conf -N
redirect_stderr = true
priority        = 600
startretries    = 10
autorestart	= true